# debian/smman.postrm
#!/bin/sh

set -e

case "$1" in
   remove|purge)
      update-rc.d -f smman remove
   ;;
esac

exit 0

